# vagrant/api.service
[Unit]
Description=Contractor API Server
After=network.target

[Service]
Type=simple
WorkingDir=/usr/local/contractor/api_server
ExecStart=/usr/local/contractor/api_server/api_server.py

[Install]
WantedBy=mutli-user.target
