# debian/cron.d
*/15 * * * * root [ -x /usr/lib/contractor/cron/genDNS ] && /usr/lib/contractor/cron/genDNS
*/5 * * * * root [ -x /usr/lib/contractor/cron/postMaster ] && /usr/lib/contractor/cron/postMaster
